# src/data/etc/cron.hourly/get-wx-alerts
#!/usr/bin/ucode
/* 
 * Program: get-wx-alerts
 * Purpose: Get the current NWS weather alerts for a specific zone and update /www/aam/wx.txt
 * SET YOUR NWS ZONE in the WX_ZONE VARIABLE BELOW
 */

const WX_ZONE = "AZZ548";

import * as fs from "fs";
import * as log from 'log';
const WX_URL = sprintf("https://api.weather.gov/alerts/active?zone=%s", WX_ZONE);

// Always clear the alerts file
f = fs.open("/www/aam/wx.txt", "w");
if(f) {
  f.write("", "w");
  f.close();
} else {
  log.syslog(log.LOG_INFO, "get-wx-alerts: cannot open wx.txt");
  exit();
}

// Test for Internet connection, exit if none
const p = fs.popen("exec /bin/ping -W1 -c1 8.8.8.8 2>/dev/null");
if (p) {
  const d = p.read("all");
  p.close();
  if (index(d, "1 packets received") == -1) {
    log.syslog(log.LOG_INFO, "get-wx-alerts: no Internet available");
    exit();
  }
}

// Query for active alert info
const wx = fs.popen(`exec /usr/bin/curl --silent --retry 3 --connect-timeout 5 --speed-time 10 --speed-limit 1000 ${WX_URL}`);
if(wx) { // NWS forecast returned
  const j = json(wx.read("all"));  //Parse the json retrieved
  wx.close();

  if(j.features[0]) { //has active alert text
    f = fs.open("/www/aam/wx.txt", "w");
    if(f) { f.write(sprintf("%s<br>%s\n", j.title, j.features[0].properties.headline)); }
    f.close();
  }
}
